# docker/supervisor/dev.conf
[include]
files = /etc/supervisor/common.conf

[program:octane]
command=php /var/www/html/artisan octane:frankenphp --watch --port=80 --caddyfile=/etc/caddy/Caddyfile
autostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
